# board/common/rootfs/etc/profile.d/alias.sh
alias cli='clish'
alias cfg='sysrepocfg -f json'
