# board/common/rootfs/etc/profile.d/convenience.sh
alias l='ls -CF'
alias la='ls -A'
alias ll='ls -alF'
alias ls='ls --color=auto'

alias less=most
alias vim='vi'
alias view='vi -R'
alias emacs='mg'
alias edit='mg'
alias editor=edit
alias sensible-editor=edit
alias sensible-pager=less
export EDITOR=/usr/bin/mg

alias ip='ip --color=auto'
alias ipb='ip -br'
alias ipaddr='ip addr'
alias iplink='ip link'
alias bridge='bridge --color=auto'
