# board/common/rootfs/etc/profile.d/update-check.sh
if [ -s /run/infix-update ]; then
    printf '\n\033[1;33m *** %s ***\033[0m\n\n' "$(cat /run/infix-update)"
fi
