# board/common/rootfs/etc/udev/rules.d/70-rename-wifi.rules
SUBSYSTEM=="net", ACTION=="add", TEST=="/sys/class/net/$name/wireless", NAME="wifi%n" RUN+="/usr/libexec/confd/wifi-detected.sh"
