# board/netconf/rootfs/etc/profile.d/cli-hint.sh
if [ -z "$CLISH" ]; then
    cat <<EOF

Run the command 'cli' for interactive OAM

EOF
fi
